# root/etc/services.d/tailscaled/run
#!/usr/bin/with-contenv bash

exec \
    /usr/local/bin/tailscaled --state=/config/tailscaled.state